# packages/canopy_runner/launchd/com.canopy.runner.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key><string>com.canopy.runner</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/bin/env</string>
    <string>python3</string>
    <string>-m</string>
    <string>canopy_runner.main</string>
    <string>--config</string>
    <string>/Users/jjackson/.canopy/runner.json</string>
  </array>
  <key>RunAtLoad</key><true/>
  <key>KeepAlive</key><true/>
  <key>StandardOutPath</key><string>/Users/jjackson/.canopy/runner.log</string>
  <key>StandardErrorPath</key><string>/Users/jjackson/.canopy/runner.log</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PYTHONPATH</key>
    <string>/Users/jjackson/emdash-projects/canopy-web/packages/canopy_runner</string>
  </dict>
</dict>
</plist>
